# scripts/macos/com.andrea.openclaw-gateway-refresh.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.andrea.openclaw-gateway-refresh</string>
  <key>RunAtLoad</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>command -v openclaw &amp;&amp; openclaw gateway restart || true</string>
  </array>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/andrea/openclaw-gateway-refresh.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/andrea/openclaw-gateway-refresh.err</string>
</dict>
</plist>
