# scripts/macos/com.andrea.openclaw-gateway-refresh.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.andrea.openclaw-gateway-refresh</string>
  <key>RunAtLoad</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>export PATH="${HOME}/.npm-global/bin:${HOME}/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:${PATH}"; if command -v openclaw &gt;/dev/null 2&gt;&amp;1; then openclaw gateway restart; else echo "openclaw not found on PATH" &gt;&amp;2; exit 127; fi</string>
  </array>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/andrea/openclaw-gateway-refresh.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/andrea/openclaw-gateway-refresh.err</string>
</dict>
</plist>
