# scripts/macos/com.andrea.andrea-cloudflared.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.andrea.andrea-cloudflared</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/cloudflared</string>
    <string>tunnel</string>
    <string>--no-autoupdate</string>
    <string>run</string>
    <string>--token</string>
    <string>__CLOUDFLARED_TUNNEL_TOKEN__</string>
  </array>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/andrea/cloudflared.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/andrea/cloudflared.err</string>
</dict>
</plist>
