# scripts/macos/com.andrea.andrea-sync.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.andrea.andrea-sync</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>15</integer>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>export PATH="${HOME}/.npm-global/bin:${HOME}/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:${PATH}"; test -f "__REPO_ROOT__/.env" &amp;&amp; source "__REPO_ROOT__/.env"; test -f "${HOME}/andrea-lockstep.env" &amp;&amp; source "${HOME}/andrea-lockstep.env"; cd "__REPO_ROOT__" &amp;&amp; exec "__PYTHON3__" scripts/andrea_sync_server.py</string>
  </array>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/andrea/andrea-sync.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/andrea/andrea-sync.err</string>
</dict>
</plist>
